# templates/alerta.service.j2
[Unit]
Description=Alerta Service

[Service]
ExecStart={{ alerta_base_dir }}/bin/uwsgi --ini /etc/uwsgi.ini

[Install]
WantedBy=multi-user.target
